# lowbattery.service
[Unit]
Description="Alert for low battery events"

[Service]
ExecStart=/usr/bin/lowbattery.sh
