# systemd/usr/lowbattery.service
[Unit]
Description="Alert for low battery events"

[Service]
ExecStart=/usr/local/bin/lowbattery.sh

[Install]
WantedBy=timer-twominutes.target
