# systemd/usr/timer-twominutes.target
[Unit]
Description=Two Minutes Timer Target
StopWhenUnneeded=yes
